# backup/Library/LaunchAgents/com.vitkutny.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
	<dict>
		<key>Label</key>
		<string>com.vitkutny.backup</string>
		<key>Program</key>
		<string>/Users/vitkutny/.local/bin/backup</string>
		<key>StandardOutPath</key>
		<string>/Users/vitkutny/Library/Logs/com.vitkutny.backup/backup.log</string>
		<key>StandardErrorPath</key>
		<string>/Users/vitkutny/Library/Logs/com.vitkutny.backup/backup.log</string>
		<key>RunAtLoad</key>
		<false/>
		<key>StartInterval</key>
		<integer>28800</integer>
	</dict>
</plist>

